# supervisord.d/crontab.ini
[program:crontab]
command=crond -n
startsecs=2
autostart=true
autorestart=true